# templates/etc/profile.d/zookeeper.sh.j2
export PATH=$PATH:{{zookeeper_parent_install_dir}}/zookeeper/bin